# services/usb-gadget-source.service
[Unit]
Description=USB Gadget External Video Source
After=usb-gadget.service usb-gadget-stream.service
Wants=usb-gadget.service usb-gadget-stream.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=-/etc/default/usb-gadget
EnvironmentFile=-/etc/default/usb-gadget-source
Environment=SOURCE_SIZE=640x480
Environment=SOURCE_FPS=30
ExecStart=/usr/local/bin/usb-gadget-source.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
